# nginx-edge-service/data/etc/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
loglevel=debug
logfile=/var/log/supervisord.log

#[program:network]
#command=/bin/bash -c "/data/bin/lvs_dr_rs.sh"
#autorestart=false

#[program:amplify-agent]
#command=/bin/bash -c "service amplify-agent start"

[program:nginx]
command=/bin/bash -c "/usr/sbin/nginx -c \"/data/etc/nginx/nginx-${KEEPALIVED_ROLE}.conf\" -g \"daemon off;\""
environment=LD_LIBRARY_PATH="/usr/local/lib:${LD_LIBRARY_PATH}"

#[program:keepalived]
#command=/bin/bash -c "/usr/local/keepalived/sbin/keepalived -n -l -D -f \"/data/etc/keepalived/keepalived-${KEEPALIVED_ROLE}.conf\""
