# api-standalone/data/etc/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
loglevel=debug
logfile=/app/data/var/log/supervisord.log

[program:api]
command=/root/anaconda3/bin/python "/app/api_standalone.py"
redirect_stderr = true
stdout_logfile=/app/data/var/log/api_standalone_stdout.log
stdout_logfile_maxbytes = 2MB  
stdout_logfile_backups = 5
; stopasgroup=true ; send stop signal to the UNIX process
; killasgroup=false ; SIGKILL the UNIX process group (def false)

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:9001        ; (ip_address:port specifier, *:port for all iface, 若的形式*:port则开放外网访问 )
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))
